# services/scripts/com.meridiona.pm-worklog-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    com.meridiona.pm-worklog-daemon — hourly Jira worklog posting daemon.

    Install via services/scripts/install-pm-worklog-daemon.sh which copies
    this file into ~/Library/LaunchAgents/ with absolute paths substituted in,
    then bootstraps it under launchd. The plist below uses {{HOME}} +
    {{REPO_ROOT}} placeholders that the install script replaces.

    Uninstall:
        services/scripts/uninstall-pm-worklog-daemon.sh
    Or manually:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.meridiona.pm-worklog-daemon.plist
        rm ~/Library/LaunchAgents/com.meridiona.pm-worklog-daemon.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.pm-worklog-daemon</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}/services</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{REPO_ROOT}}/services/.venv/bin/python</string>
        <string>-m</string>
        <string>agents.pm_worklog_update.jira_worklog_daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MERIDIAN_DB</key>
        <string>{{HOME}}/.meridian/meridian.db</string>
        <key>MERIDIAN_HOME</key>
        <string>{{HOME}}/.meridian</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <!-- Worklog posting tunables. PM_WORKLOG_INTERVAL_HOURS controls how
             often the daemon wakes to post worklogs. PM_WORKLOG_LOOKBACK_HOURS
             controls the session window inspected on each run. -->
        <key>PM_WORKLOG_INTERVAL_HOURS</key>
        <string>1</string>
        <key>PM_WORKLOG_LOOKBACK_HOURS</key>
        <string>1</string>
        <!-- Jira credentials — substituted from services/.env by the install
             script. Required for the Jira REST API worklog endpoints. -->
        <key>JIRA_URL</key>
        <string>{{JIRA_URL}}</string>
        <key>JIRA_EMAIL</key>
        <string>{{JIRA_EMAIL}}</string>
        <key>JIRA_API_TOKEN</key>
        <string>{{JIRA_API_TOKEN}}</string>
        <!-- OTLP auth for OpenObserve — must be set here (not just .env)
             because observability.py initialises before dotenv is loaded.
             The install script reads MERIDIAN_OO_AUTH from services/.env and
             substitutes the {{MERIDIAN_OO_AUTH}} placeholder below. -->
        <key>MERIDIAN_OO_AUTH</key>
        <string>{{MERIDIAN_OO_AUTH}}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/pm-worklog-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/pm-worklog-daemon.log</string>

    <!-- RunAtLoad true: catches up immediately on login, then the daemon
         manages its own interval sleep. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
