# services/scripts/com.meridiona.ax-sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.ax-sidecar</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PYTHON}}</string>
        <string>{{REPO_ROOT}}/services/ax_sidecar.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SCREENPIPE_DB</key>
        <string>{{HOME}}/.screenpipe/db.sqlite</string>
        <key>AX_SIDECAR_INTERVAL</key>
        <string>3</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/ax-sidecar.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/ax-sidecar.log</string>

    <key>ThrottleInterval</key>
    <integer>2</integer>
</dict>
</plist>
