# services/scripts/com.meridiona.mlx-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    com.meridiona.mlx-server — persistent MLX inference server.

    Loads the Qwen3.5-9B model once at startup and serves POST /classify_sessions
    to the Rust daemon. KeepAlive keeps it running across crashes; launchd
    restarts it automatically.

    Install:
        services/scripts/install-mlx-server-daemon.sh
    Uninstall:
        services/scripts/uninstall-mlx-server-daemon.sh
    Status / logs:
        launchctl print gui/$(id -u)/com.meridiona.mlx-server
        tail -f ~/.meridian/logs/mlx-server.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.mlx-server</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}/services</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{REPO_ROOT}}/services/.venv/bin/meridian-server</string>
        <string>--backend</string>
        <string>mlx</string>
        <string>--port</string>
        <string>{{MLX_SERVER_PORT}}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MERIDIAN_DB</key>
        <string>{{HOME}}/.meridian/meridian.db</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>MLX_SERVER_PORT</key>
        <string>{{MLX_SERVER_PORT}}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/mlx-server.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/mlx-server.log</string>

    <!-- Start immediately when loaded and after every crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Wait 5 s before restarting after a crash (avoids tight crash loops). -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <!-- ProcessType omitted intentionally: Background restricts file access to ~/Documents
         on macOS 15+ (sandbox blocks pyvenv.cfg reads). Standard QoS is correct here. -->
</dict>
</plist>
