# services/scripts/com.meridiona.tagger-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    com.meridiona.tagger-daemon — long-running 3-stage tagger.

    Install via services/scripts/install-tagger-daemon.sh which copies this
    file into ~/Library/LaunchAgents/ with absolute paths substituted in,
    then bootstraps it under launchd. The plist below uses {{HOME}} +
    {{REPO_ROOT}} placeholders that the install script replaces.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.tagger-daemon</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}/services</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{REPO_ROOT}}/services/.venv/bin/python</string>
        <string>-m</string>
        <string>agents.tagger_daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MERIDIAN_DB</key>
        <string>{{HOME}}/.meridian/meridian.db</string>
        <key>MERIDIAN_HOME</key>
        <string>{{HOME}}/.meridian</string>
        <key>LOG_LEVEL</key>
        <string>INFO</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <!-- ONLY_TODAY=1 — daemon only tags sessions started today (local time).
             Set ONLY_TODAY=0 here to force backfill of historical sessions. -->
        <key>ONLY_TODAY</key>
        <string>1</string>
        <!-- Tick is also configurable via the agents/config.py defaults; this
             override keeps it close to the supervisor for easy tuning. -->
        <key>TAGGER_TICK_SECS</key>
        <string>7</string>
        <!-- Stage 3 (LLM) credentials are read from .env files inside config.py;
             the launchd context still loads them since we set MERIDIAN_HOME and
             agents/config.py walks ~/.hermes/.env automatically. -->
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/tagger-daemon.out</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/tagger-daemon.err</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Don't loop a misconfigured daemon at 100% CPU. 10s minimum between
         restarts matches the launchd recommendation. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <!-- Background scheduling priority — never compete with foreground UI work. -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
