# scripts/com.meridiona.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    com.meridiona.daemon — main Meridian Rust daemon. Runs continuously, polls screenpipe.

    Install via scripts/install-daemon.sh which copies this file into
    ~/Library/LaunchAgents/ with absolute paths substituted in, then
    bootstraps it under launchd. The plist below uses {{HOME}},
    {{REPO_ROOT}}, {{DAEMON_BIN}}, {{MERIDIAN_OO_AUTH}}, and
    {{MERIDIAN_OTLP_ENDPOINT}} placeholders that the install script replaces.

    Uninstall:
        scripts/uninstall-daemon.sh
    Or manually:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.meridiona.daemon.plist
        rm ~/Library/LaunchAgents/com.meridiona.daemon.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.daemon</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{DAEMON_BIN}}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd gives agents a minimal PATH (/usr/bin:/bin:/usr/sbin:/sbin).
             The coding-agent summariser shells out to `claude` (~/.local/bin) and
             `codex` (/opt/homebrew/bin); without them on PATH every primary
             summary silently falls back to MLX. Keep ~/.local/bin and Homebrew
             ahead of the system defaults. -->
        <key>PATH</key>
        <string>{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SCREENPIPE_DB</key>
        <string>{{HOME}}/.screenpipe/db.sqlite</string>
        <key>MERIDIAN_DB</key>
        <string>{{HOME}}/.meridian/meridian.db</string>
        <key>POLL_INTERVAL_SECS</key>
        <string>60</string>
        <key>RUST_LOG</key>
        <string>meridian=info</string>
        <!-- OTLP auth for OpenObserve — must be set here (not just .env)
             because the Rust binary reads env vars before dotenv is loaded.
             The install script reads MERIDIAN_OO_AUTH from the repo-root .env and
             substitutes the {{MERIDIAN_OO_AUTH}} placeholder below. -->
        <key>MERIDIAN_OO_AUTH</key>
        <string>{{MERIDIAN_OO_AUTH}}</string>
        <key>MERIDIAN_OTLP_ENDPOINT</key>
        <string>{{MERIDIAN_OTLP_ENDPOINT}}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/daemon-error.log</string>

    <!-- RunAtLoad true: the Rust daemon must always be running — it fires
         immediately on login and is kept alive by launchd. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
