# scripts/com.meridiona.ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    com.meridiona.ui — runs the Next.js dashboard (built by ./install.sh) as
    a background launchd LaunchAgent. Serves on http://localhost:3939.

    Install via scripts/install-ui-daemon.sh which copies this file into
    ~/Library/LaunchAgents/ with absolute paths substituted in, then
    bootstraps it under launchd. The plist below uses {{HOME}}, {{REPO_ROOT}}
    and {{NPM_BIN}} placeholders that the install script replaces.

    Uninstall:
        scripts/uninstall-ui-daemon.sh
    Or manually:
        launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.meridiona.ui.plist
        rm ~/Library/LaunchAgents/com.meridiona.ui.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meridiona.ui</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}/ui</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{NPM_BIN}}</string>
        <string>run</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>MERIDIAN_DB</key>
        <string>{{HOME}}/.meridian/meridian.db</string>
        <key>PORT</key>
        <string>3939</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.meridian/logs/ui.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.meridian/logs/ui-error.log</string>

    <!-- RunAtLoad true: dashboard should always be available after login. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
